# dot.config/systemd/user/wallpaper.timer
[Unit]
Description=Switch Wallpapers Every 15 minutes

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=30min

[Install]
WantedBy=wm.target
