# dot.config/systemd/user/wm.target
[Unit]
Description=Window manager
AllowIsolate=yes
